# deploy/supervisor/amoeba_prd.ini
[program:amoeba_prd]
command=/data/prd/amoeba/env/bin/gunicorn -c gunicorn.conf.py -p gunicorn-amoeba.pid  config.wsgi
directory=/data/prd/amoeba/src
user=pyer
autostart=true
autorestart=true
startretries=5
stdout_logfile=/data/prd/amoeba/logs/supervisor_stdout.log
stderr_logfile=/data/prd/amoeba/logs/supervisor_stderr.log
environment=DEBUG="False"